# debian/postinst
#!/bin/sh
# postinst script for ganglia-modules-linux
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        INIT_GMOND=/etc/init.d/ganglia-monitor
        if which invoke-rc.d >/dev/null 2>&1; then
    		invoke-rc.d ganglia-monitor status >/dev/null 2>&1 && \
     			invoke-rc.d ganglia-monitor restart
     	else
        	[ -x ${INIT_GMOND} ] && \
                	${INIT_GMOND} status >/dev/null 2>&1 && \
                	${INIT_GMOND} restart
     	fi
        # call `true' to conceal the return status of the gmond restart
        # as we don't want an unrelated error in gmond to suggest
        # that the dpkg operation has completely failed
        true
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
